# deploy/data-ops/darlink-dashboard.service
[Unit]
Description=Darlink — Streamlit data dashboard (long-running)
After=network.target
ConditionPathExists=/opt/darlink/data/darlink.sqlite

[Service]
Type=simple
User=darlink
Group=darlink
WorkingDirectory=/opt/darlink
ExecStart=/opt/darlink/.venv/bin/streamlit run /opt/darlink/scripts/dashboard.py \
    --server.headless true \
    --server.port 8501 \
    --server.address 0.0.0.0 \
    --browser.gatherUsageStats false
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/darlink/scripts/logs/dashboard.log
StandardError=append:/opt/darlink/scripts/logs/dashboard.err

# 简单的资源限制 + 安全加固
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
